# external/external_v2.0/etcd_installer/etcd_clone.service
[Unit]
Description=etcd key-value store
Documentation=https://github.com/etcd-io/etcd
After=network.target

[Service]
User=root
Type=notify
Environment=ETCD_DATA_DIR=/var/lib/etcd_clone
Environment=ETCD_NAME=%m
ExecStart=/usr/local/bin/etcd \
  --cert-file=/etc/pki/etcd-ca/certs/etcd0.example.com.crt \
  --key-file=/etc/pki/etcd-ca/private/etcd0.example.com.key \
  --trusted-ca-file=/etc/pki/etcd-ca/certs/ca.crt \
  --listen-client-urls=https://<IPADDRESS>:2369,https://127.0.0.1:2369 \
  --advertise-client-urls=https://<IPADDRESS>:2369 \
  --listen-peer-urls=http://<IPADDRESS>:2370 \
  --initial-advertise-peer-urls=http://<IPADDRESS>:2370 \
  --initial-cluster=default=http://<IPADDRESS>:2370 \
  --data-dir=/var/lib/etcd_clone
Restart=always
RestartSec=10s
LimitNOFILE=40000

[Install]
WantedBy=multi-user.target
